# syscfg/etc/crontabs/root
# Cron Table
#
#
04 03 * * * /etc/init.d/rc.daily		# Daily service script

